# precise-cassandra-standalone/confs/supervisor/conf.d/cassandra.conf
[program:cassandra]
command=/bin/sh /etc/init.d/cassandra start