# precise/precise-elasticsearch/confs/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch]
command=/bin/sh /etc/init.d/elasticsearch start