# precise-hbase-standalone/confs/supervisor/conf.d/hbase.conf
[program:hbase]
command=/bin/sh /usr/local/hbase/bin/start-hbase.sh