# precise/precise-storm-pseudo/confs/supervisor/conf.d/storm_supervisor.conf
[program:storm_supervisor]
command=/usr/local/storm/bin/storm supervisor