# precise/precise-storm-pseudo/confs/supervisor/conf.d/storm_nimbus.conf
[program:storm_nimbus]
command=/usr/local/storm/bin/storm nimbus