# precise/precise-storm-pseudo/confs/supervisor/conf.d/storm_ui.conf
[program:storm_ui]
command=/usr/local/storm/bin/storm ui