# precise/precise-storm-pseudo/confs/supervisor/conf.d/zookeeper.conf
[program:zookeeper]
command=/bin/sh /usr/local/zookeeper/bin/zkServer.sh start