# precise/precise-hbase-standalone/confs/supervisor/conf.d/hbase-thrift.conf
[program:hbase-thrift]
command=/bin/sh /usr/local/hbase/bin/hbase-daemon.sh start thrift