# precise/precise-hbase-standalone/confs/supervisor/conf.d/hbase-rest.conf
[program:hbase-rest]
command=/bin/sh /usr/local/hbase/bin/hbase-daemon.sh start rest