# precise-cdh4-mrv2/confs/supervisor/conf.d/hadoop-yarn-resourcemanager.conf
[program:hadoop-yarn-resourcemanager]
command=/etc/init.d/hadoop-yarn-resourcemanager start