# precise/precise-hbase-mrv2-pseudo/confs/supervisor/conf.d/hadoop-yarn-nodemanager.conf
[program:hadoop-yarn-nodemanager]
command=/etc/init.d/hadoop-yarn-nodemanager start