# precise/precise-hbase-mrv2-pseudo/confs/supervisor/conf.d/hbase-thrift.conf
[program:hbase-thrift]
command=/bin/sh /etc/init.d/hbase-thrift start