# precise/precise-hbase-mrv2-pseudo/confs/supervisor/conf.d/hbase-regionserver.conf
[program:hbase-regionserver]
command=/bin/sh /etc/init.d/hbase-regionserver start