# precise/precise-hbase-mrv2-pseudo/confs/supervisor/conf.d/hbase-rest.conf
[program:hbase-rest]
command=/bin/sh /etc/init.d/hbase-rest start