# precise/precise-hbase-mrv2-pseudo/confs/supervisor/conf.d/zookeeper-server.conf
[program:zookeeper-server]
command=/bin/sh /etc/init.d/zookeeper-server start