# precise/precise-cdh4-mrv2/confs/supervisor/conf.d/hadoop-hdfs-namenode.conf
[program:hadoop-hdfs-namenode]
command=/etc/init.d/hadoop-hdfs-namenode start