# precise/precise-cdh4-mrv2/confs/supervisor/conf.d/hadoop-mapreduce-historyserver.conf
[program:hadoop-mapreduce-historyserver]
command=/etc/init.d/hadoop-mapreduce-historyserver start