# precise-hbase-titan/confs/supervisor/conf.d/rexster-server.conf
[program:rexster-server]
command=/bin/sh /usr/local/rexster/bin/rexster-service.sh start