# precise-cdh4-mrv1/confs/supervisor/conf.d/hadoop-hdfs-secondarynamenode.conf
[program:hadoop-hdfs-secondarynamenode]
command=/etc/init.d/hadoop-hdfs-secondarynamenode start