# precise-hbase-mrv2-pseudo/confs/supervisor/conf.d/hbase-master.conf
[program:hbase-master]
command=/bin/sh /etc/init.d/hbase-master start