# precise/precise-cdh4-mrv1/confs/supervisor/conf.d/hadoop-0_20-mapreduce-jobtracker.conf
[program:hadoop-0_20-mapreduce-jobtracker]
command=/etc/init.d/hadoop-0.20-mapreduce-jobtracker start