# precise/precise-cdh4-mrv1/confs/supervisor/conf.d/hadoop-0_20-mapreduce-tasktracker.conf
[program:hadoop-0_20-mapreduce-tasktracker]
command=/etc/init.d/hadoop-0.20-mapreduce-tasktracker start