# precise-cdh4-hbase-standalone/confs/supervisor/conf.d/hbase-master.conf
[program:hbase-master]
command=/etc/init.d/hbase-master start