# precise-cdh4-hbase-standalone/confs/supervisor/conf.d/hbase-rest.conf
[program:hbase-rest]
command=/etc/init.d/hbase-rest start